# ops/com.eshen.agent-ready-outreach.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.eshen.agent-ready-outreach</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/bigo/Desktop/agent-ready/scripts/run-outreach-operator.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/bigo/Desktop/agent-ready</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/bigo/Desktop/agent-ready/.agent-ready/outreach/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bigo/Desktop/agent-ready/.agent-ready/outreach/launchd.err.log</string>
</dict>
</plist>
